# upstart-nginx.conf
description "nginx-rp container"
author "chad.neal@jeppesen.com"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a nginx-rp
end script
